# app/babycare/supervisor/babycare-monitor.conf
[program:babycare-monitor]
command=python3 /path/to/myproject/camera_local_file.py
directory=/path/to/myproject/
user=root
autostart=true
autorestart=true
redirect_stderr=true